# nginx-php/service/uwsgi/run
#!/bin/sh
exec 2>&1
exec /usr/bin/uwsgi --ini /etc/uwsgi.ini
